# bubbletube-update.service
[Unit]
Description=BubbleTube Update Service

[Service]
ExecStart=/home/pi/bubbletube/update-git
Type=oneshot

[Install]
WantedBy=multi-user.target

